# autostart/autostart_gropo.desktop
[Desktop Entry]

Exec=sudo bash /home/aronn/GroPro_onstart.sh